# bin/crontab/archive-tmp
#!/bin/sh
DATE=$(date +%Y-%m-%d)
CURRENT_TMP=~/tmp
ARCHIVE=~/archive-tmp/

mkdir -p $ARCHIVE
mv -n "$CURRENT_TMP" "$ARCHIVE/$DATE"
mkdir "$CURRENT_TMP"
